# desktopbackground.service
[Unit]
Description="Desktop Background Randomizer"

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'desktop-bg-rand.sh | set-background'
