# desktopbackground.timer
[Unit]
Description="Desktop Background Randomization Timer"

[Timer]
OnActiveSec=0
OnUnitInactiveSec=0

[Install]
WantedBy=timers.target
