# desktopbackground.timer
[Unit]
Description="Desktop Background Randomization Timer"

[Timer]
OnActiveSec=0
OnUnitActiveSec=120

[Install]
WantedBy=timers.target
