# desktopbackground.service
[Unit]
Description=Desktop Background Randomizer

[Service]
Type=simple
Environment=SPEED=20
ExecStart=/bin/sh -c '/usr/bin/desktop-bg-rand.sh --print-rating --multiply-rating ${SPEED} | /usr/bin/set-background --wait-after'
Restart=always

[Install]
WantedBy=default.target
